# devtools/builder/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

* * * * * root cd /code && python ./schedule.py > /proc/1/fd/1 2>/proc/1/fd/2